# conf/supervisor/jwschain/JobReportTaxReceipt.ini
[program:laravel-JobReportTaxReceipt]
process_name=%(program_name)s_%(process_num)02d
directory=/www/wwwroot/
command=php artisan queue:work redis --queue=JobReportTaxReceipt ; 被监控的进程启动命令
autostart=true ; 随着supervisord的启动而启动
autorestart=true ; 自动重启.
user=www ; supervisord以root运行
numprocs=1 ; 启动几个进程
redirect_stderr=true ; 重定向stderr到stdout
stdout_logfile=/tmp/JobReportTaxReceipt.log ; 进程的标准输出写入文件中
